# rtl4332mqtt/supervisord.conf
[supervisord]
user= root
loglevel= warn

[supervisorctl]

[program:rtl2mqtt]
command= /rtl2mqtt.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autorestart=true
